# cockpit.container
[Container]
Image=quay.io/cockpit/ws
ContainerName=cockpit
Environment=COCKPIT_SSH_KEY_PATH=/run/secrets/COCKPIT_SSH
Volume=cockpit.volume:/etc/cockpit
Network=cockpit.network
PodmanArgs=--secret COCKPIT_SSH,type=mount,mode=0400
Label=io.containers.autoupdate=registry

[Service]
Restart=always

[Install]
WantedBy=default.target